# resources/lib/systemd/system/lmd-update.service
[Unit]
Description=lmd Linux automatic update

[Service]
Type=oneshot
ExecStart=/lib/lmd/auto-update
